# scripts/ss_proxy_setup/roles/client/templates/shadowsocks.client.service.j2
[Unit]
Description=shadowsocks proxy server

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/bin/ss-local -c {{ proxy_client_config }}

[Install]
WantedBy=multi-user.target
