# scripts/ss_proxy_setup/roles/server/templates/shadowsocks.server.service.j2
[Unit]
Description=shadowsocks proxy server

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/bin/ss-server -c {{ proxy_server_config }}

[Install]
WantedBy=multi-user.target
